# deploy/tmo-ai.service
[Unit]
Description=TMO-AI FastAPI service
After=network.target

[Service]
# Run as an unprivileged system user created for the app (adjust as needed)
User=tmoai
Group=tmoai

# Working directory containing backend (where tmo_ai.db will be created by default)
WorkingDirectory=/opt/tmo-ai/backend

# Ensure the virtualenv's bin is used for ExecStart and any child processes
Environment="PATH=/opt/tmo-ai/venv/bin"

# Cache directories for Hugging Face / Transformers to avoid writes to /home/tmoai
# Create and chown `/var/cache/tmo-ai` on the server and make it writable by the `tmoai` user.
Environment="TRANSFORMERS_CACHE=/var/cache/tmo-ai/huggingface"
Environment="XDG_CACHE_HOME=/var/cache/tmo-ai"
Environment="HOME=/opt/tmo-ai"
Environment="HF_HOME=/var/cache/tmo-ai/huggingface"
Environment="HUGGINGFACE_HUB_CACHE=/var/cache/tmo-ai/huggingface"
Environment="HF_DATASETS_CACHE=/var/cache/tmo-ai/huggingface"

# Example: set database path to a writable location (optional)
# Environment="TMO_AI_DB_PATH=/var/lib/tmo-ai/tmo_ai.db"

# Command to start the app. Adjust host/port/workers as needed.
ExecStart=/opt/tmo-ai/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4

Restart=on-failure
RestartSec=5

# Increase file descriptors if needed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
